# data/systemd/redshift.service.in
[Unit]
Description=Redshift display colour temperature adjustment
Documentation=http://jonls.dk/redshift/
After=display-manager.service

[Service]
ExecStart=@bindir@/redshift
Restart=always

[Install]
WantedBy=default.target
