# data/systemd/redshift-gtk.service.in
[Unit]
Description=Redshift display colour temperature adjustment (GUI)
Documentation=http://jonls.dk/redshift/
After=graphical-session.target

[Service]
ExecStart=@bindir@/redshift-gtk
Restart=always

[Install]
WantedBy=default.target
