# data/systemd/redshift.service.in
[Unit]
Description=Redshift display colour temperature adjustment
Documentation=http://jonls.dk/redshift/
After=graphical-session.target

[Service]
ExecStart=@bindir@/redshift
Restart=always
RestartSec=30

[Install]
WantedBy=graphical-session.target
