# data/upstart/redshift.conf
description "Redshift display colour temperature adjustment"
author "Cameron Norman <camerontnorman@gmail.com>"

start on startup
stop on session-end

exec redshift
